# scheduler/business-scheduler.service
[Unit]
Description=Business Workflow Scheduler
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/business-scheduler --app-root %h/workspaces/jambu/business-workflow/app
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
